# deploy/autostart/lmpanel.desktop
[Desktop Entry]
Type=Application
Name=Linux Monitor Panel
Comment=Desktop information overlay
Exec=/bin/sh -lc 'systemctl --user import-environment DISPLAY WAYLAND_DISPLAY XAUTHORITY DBUS_SESSION_BUS_ADDRESS; systemctl --user start lmpanel.service'
Terminal=false
StartupNotify=false
X-GNOME-Autostart-enabled=true
OnlyShowIn=GNOME;XFCE;MATE;Cinnamon;Unity;Budgie;
